# daemons/fluidsynth.service
[Unit]
Description=FluidSynth Daemon
After=sound.target

[Service]
ExecStart=/usr/bin/fluidsynth -is -o synth.midi-bank-select=xg -a alsa -g 1 '/home/pi/DecembreSF.sf2'

[Install]
WantedBy=multi-user.target
